# dotfiles/.config/systemd/user/autostart-solaar.service
[Unit]
Description=Auto-start Solaar when any Bluetooth device connects
After=bluetooth.target

[Service]
ExecStart=/usr/local/bin/autostart-solaar.sh
Restart=on-failure

[Install]
WantedBy=default.target
